# conf/github-contributions-process.service
[Unit]
Description=GitHub Contributions Process Task

[Service]
ExecStart=/srv/github-contributions/conf/process-job
