# conf/github-contributions.conf
# Upstart script for GHC Server
# Belongs in /etc/init

description "GitHub Contributions Server"
author "Liam Bowen"
start on runlevel [2345]
stop on runlevel [06]

setuid ghc
setgid ghc

respawn

env GHC_ENV=production

exec /srv/bin/ghc-app
