# conf/github-contributions-process.timer
[Unit]
Description=Run GHC archive processor hourly

[Timer]
Persistent=true
OnCalendar=hourly

[Install]
WantedBy=timers.target